# provisioning/roles/golang/files/etc/systemd/system/isucon-app.golang.service
[Unit]
Description = isucon7final main application in golang

[Service]
WorkingDirectory=/home/isucon/isucon-all/webapp/go/
EnvironmentFile=/home/isucon/env.sh
ExecStart = /home/isucon/isucon-all/webapp/go/app

Restart   = always
Type      = simple
User      = isucon
Group     = isucon

[Install]
WantedBy = multi-user.target
